# files/root/etc/systemd/system/last-pipeline2.service
[Unit]
Description=LAST pipeline service 2 (of 2)

[Service]
User=ocs
WorkingDirectory=/home/ocs/matlab
ExecStart=/usr/local/share/last-tool/bin/last-pipeline start 2
ExecStop=/usr/local/share/last-tool/bin/last-pipeline stop 2
Restart=always
Environment="SYSTEMD=1"
WatchdogSec=1800
NotifyAccess=all
Type=forking

[Install]
WantedBy=multi-user.target
